# etc/cron.d/1duckdns
SHELL=/bin/bash
PATH=/usr/bin
MAILTO=icewootra
* * * * * root duckdns
