# home/icewootra/.config/systemd/user/alarm.timer
[Unit]
Description=Alarm at 8:00

[Timer]
Unit=spotify.service
OnCalendar=Mon-Sun *-*-* 7:00

[Install]
WantedBy=timers.target
