# home/icewootra/.config/systemd/user/dbus-org.jackaudio.service
[Unit]
Description=JACK Audio system DBUS activation service.

[Service]
Type=dbus
#ExecStart=/usr/bin/jackdbus auto
BusName=org.jackaudio.service
