# home/icewootra/.config/systemd/user/spotify.service
[Unit]
Description=Play spotify!

[Service]
ExecStart=/home/icewootra/spotify.sh play
