# etc/systemd/system/netctl@ssh_tap.service
.include /usr/lib/systemd/system/netctl@.service

[Unit]
Description=tuntap connection for ssh
