# etc/systemd/system/netctl@openvpn_tap.service
.include /usr/lib/systemd/system/netctl@.service

[Unit]
Description=tuntap connection - openvpn
